# os/sysconfig/bloom-pi-agent.service
[Unit]
Description=Bloom Pi Agent (persistent Matrix listener)
After=network-online.target bloom-matrix.service
Wants=network-online.target

[Service]
Type=simple
# Pi needs a PTY for its TUI — script provides one in headless mode
ExecStartPre=/usr/local/bin/bloom-greeting.sh
ExecStart=/usr/bin/script -qfc "pi --continue" /dev/null
Environment=HOME=%h
Environment=BLOOM_DIR=%h/Bloom
Environment=PI_SESSION=daemon
Environment=TERM=dumb
Restart=on-failure
RestartSec=30
# Don't start until the user has completed first-boot setup
ConditionPathExists=%h/.bloom/.setup-complete

[Install]
WantedBy=default.target
