# services/cinny/quadlet/bloom-cinny.container
[Unit]
Description=Bloom Cinny - browser-based Matrix web client
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/cinnyapp/cinny:v4.10.2
ContainerName=bloom-cinny

PublishPort=8081:80

Volume=%h/.config/bloom/cinny/config.json:/app/config.json:ro
Volume=%h/.config/bloom/cinny/.well-known/matrix/client:/app/.well-known/matrix/client:ro
Volume=%h/.config/bloom/cinny/nginx.conf:/etc/nginx/conf.d/default.conf:ro

PodmanArgs=--memory=192m
PodmanArgs=--security-opt label=disable
NoNewPrivileges=true
LogDriver=journald

[Service]
Restart=on-failure
RestartSec=10
TimeoutStartSec=120

[Install]
WantedBy=default.target
